# kor-webdav.path
[Unit]
Description=Description=import data from a public webdav to ConedaKOR (path trigger)

[Path]
PathExists=/home/kor/webdav/new/start.txt
Unit=kor-webdav.service

[Install]
WantedBy=multi-user.target
